# init/bootstrap.service
[Unit]
Description=Node bootstrapping

[Service]
Type=oneshot
User=cluster
ExecStart=+/var/lib/cluster/init/bootstrap.sh
